# systemd/vps-ml-outcome-label-forwarder-loop.service
[Unit]
Description=AI Trading ML Outcome Label Forwarder Loop
After=docker.service vps-ml-ou-meanrev-dataset-join-loop.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/ai-trading-vps-bot
ExecStart=/bin/bash -lc 'while true; do docker exec ai-trading-bot bash -lc "cd /app && python scripts/ml_outcome_label_forwarder_v1.py"; sleep 600; done'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
